# launchagents/com.projectmap.restore.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.projectmap.restore</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>INSTALL_PATH_PLACEHOLDER/scripts/restore-layout.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>HOME_PLACEHOLDER/Library/Logs/project-map-restore.log</string>
    <key>StandardErrorPath</key>
    <string>HOME_PLACEHOLDER/Library/Logs/project-map-restore.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
